# Linux/renew_ioi.service
[Unit]
Description=Renews the home directory of the ioi user

[Service]
Type=oneshot
ExecStart=/home/sysoperator/bin/ioi_reset.sh
